# netlensscan.service
[Unit]
Description=NetLens API Service
After=network.target mongod.service mongodb.service
Wants=mongod.service mongodb.service

[Service]
Type=simple
User=netlens
WorkingDirectory=/opt/netlens
EnvironmentFile=-/opt/netlens/config.env
Environment="NODE_ENV=production"
ExecStart=/usr/bin/node /opt/netlens/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
