# netlens.service
[Unit]
Description=NetLens Network Discovery Service
After=network.target mongod.service mongodb.service
Wants=mongod.service mongodb.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netlens
EnvironmentFile=-/opt/netlens/config.env
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/opt/netlens/venv/bin/python /opt/netlens/scanner_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
