# netscanner.service
[Unit]
Description=NetLens Network Discovery Service
After=network.target mongodb.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netscanner
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /opt/netscanner/scanner_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
